# cns-unified/services/com.claude.cns-system.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Identity -->
    <key>Label</key>
    <string>com.claude.cns-system</string>
    
    <!-- Program Configuration -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/terryli/.claude/cns-unified/bin/cns-local-hub-launchd.py</string>
    </array>
    
    <!-- Working Directory -->
    <key>WorkingDirectory</key>
    <string>/Users/terryli/.claude/cns-unified</string>
    
    <!-- Execution Settings -->
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <!-- Restart if process exits -->
        <key>SuccessfulExit</key>
        <false/>
        <!-- Don't restart if process crashes too quickly -->
        <key>Crashed</key>
        <true/>
        <!-- Check if process is still alive every 30 seconds -->
        <key>PathState</key>
        <dict>
            <key>/tmp/cns-system-manager.pid</key>
            <true/>
        </dict>
    </dict>
    
    <!-- Logging Configuration -->
    <key>StandardOutPath</key>
    <string>/Users/terryli/.claude/cns-unified/logs/cns-launchd.out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/terryli/.claude/cns-unified/logs/cns-launchd.error.log</string>
    
    <!-- Process Settings -->
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>AbandonProcessGroup</key>
    <true/>
    
    <!-- Resource Limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>2048</integer>
        <key>NumberOfProcesses</key>
        <integer>200</integer>
    </dict>
    
    <!-- Environment Variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/terryli/.local/bin</string>
        <key>HOME</key>
        <string>/Users/terryli</string>
        <key>USER</key>
        <string>terryli</string>
        <key>SHELL</key>
        <string>/bin/zsh</string>
        <key>CNS_LOG_LEVEL</key>
        <string>INFO</string>
    </dict>
    
    <!-- Network Requirements -->
    <key>ServiceDescription</key>
    <string>CNS Remote Alert System - Maintains CNS hub and SSH tunnel for remote notifications</string>
    
    <!-- Throttle Settings -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <!-- Session Type -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
    
    <!-- Exit Timeout -->
    <key>ExitTimeOut</key>
    <integer>30</integer>
</dict>
</plist>