# automation/lychee/runtime/bot/com.terryli.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Identity -->
    <key>Label</key>
    <string>com.terryli.telegram-bot</string>

    <!-- Program Configuration -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/doppler</string>
        <string>run</string>
        <string>--project</string>
        <string>claude-config</string>
        <string>--config</string>
        <string>dev</string>
        <string>--</string>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>/Users/terryli/.claude/automation/lychee/runtime/bot/multi-workspace-bot.py</string>
    </array>

    <!-- Working Directory -->
    <key>WorkingDirectory</key>
    <string>/Users/terryli/.claude/automation/lychee/runtime/bot</string>

    <!-- Auto-start and Restart Configuration -->
    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <!-- Restart if process crashes -->
        <key>SuccessfulExit</key>
        <false/>
        <!-- Restart after 10s if process exits -->
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Restart throttling (prevent rapid restart loops) -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/telegram-bot.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/telegram-bot-error.log</string>

    <!-- Environment Variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/terryli</string>
    </dict>

    <!-- Process Priority (optional, lower = higher priority) -->
    <key>Nice</key>
    <integer>0</integer>

    <!-- Resource Limits (optional) -->
    <key>SoftResourceLimits</key>
    <dict>
        <!-- Max 1GB memory -->
        <key>MemoryLimit</key>
        <integer>1073741824</integer>
    </dict>

    <!-- Exit Timeout -->
    <key>ExitTimeOut</key>
    <integer>5</integer>
</dict>
</plist>
