# automation/lychee/runtime/bot/com.terryli.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Identity -->
    <key>Label</key>
    <string>com.terryli.telegram-bot</string>

    <!-- Program Configuration: Run watchexec (not bot directly!) -->
    <!-- launchd supervises watchexec, watchexec supervises bot -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/terryli/.claude/automation/lychee/runtime/bot/run-bot-prod-watchexec.sh</string>
    </array>

    <!-- Working Directory -->
    <key>WorkingDirectory</key>
    <string>/Users/terryli/.claude/automation/lychee/runtime/bot</string>

    <!-- Auto-start and Restart Configuration -->
    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <!-- Always keep watchexec running -->
        <key>SuccessfulExit</key>
        <false/>
        <!-- Restart if watchexec crashes -->
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Restart throttling (prevent rapid restart loops) -->
    <!-- If watchexec crashes 3+ times in 10s, something is seriously wrong -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/terryli/.claude/automation/lychee/logs/telegram-bot-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/terryli/.claude/automation/lychee/logs/telegram-bot-launchd-error.log</string>

    <!-- Environment Variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/terryli</string>
        <!-- Production mode: disable idle timeout -->
        <key>BOT_IDLE_TIMEOUT</key>
        <string>0</string>
        <!-- Alert mode: enable health monitoring -->
        <key>BOT_HEALTH_MONITORING</key>
        <string>1</string>
    </dict>

    <!-- Process Priority (normal priority for background service) -->
    <key>Nice</key>
    <integer>0</integer>

    <!-- Resource Limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <!-- Max 1GB memory for entire watchexec tree -->
        <key>MemoryLimit</key>
        <integer>1073741824</integer>
    </dict>

    <!-- Exit Timeout (give watchexec time to shutdown gracefully) -->
    <key>ExitTimeOut</key>
    <integer>10</integer>

    <!-- Run after network is available (bot needs internet) -->
    <key>LaunchOnlyOnce</key>
    <false/>
</dict>
</plist>
